# order-service/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:php-server]
command=php -S 0.0.0.0:8000 -t /var/www/html/public
directory=/var/www/html
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-server.err.log
stdout_logfile=/var/log/supervisor/php-server.out.log
user=www-data
priority=10

[program:saga-consumer]
command=php /var/www/html/artisan saga:consume
directory=/var/www/html
autostart=true
autorestart=true
startsecs=10
startretries=5
stderr_logfile=/var/log/supervisor/saga-consumer.err.log
stdout_logfile=/var/log/supervisor/saga-consumer.out.log
user=www-data
priority=20
